# crontab
56 02-23/3 * * * /root/acri-aserv/start_container.sh ag001 172.16.9.1 >> /var/log/aserv.log 2>&1
