# crontab
* * * * * /root/acri-aserv/start.sh >> /var/log/aserv.log 2>&1
