# docker/files/supervisor/entry.conf
[program:entry]
command=/usr/local/bin/docker-entrypoint.sh
priority=600
startsecs=0
autorestart=false

