# docker/files/supervisor/nis.conf
[program:rpcbind]
command=/sbin/rpcbind -wf
priority=200

[program:ypbind]
command=ypbind -n
priority=210
