# docker/files/supervisor/nis.conf
[program:rpcbind]
command=/sbin/rpcbind -w
priority=200

[program:ypbind]
command=ypbind -broadcast
priority=210
