# docker/files/supervisor/xbutler.conf
[program:xbutler]
command=/etc/init.d/xbutler start
priority=500
