# docker/files/supervisor/xrm.conf
[program:xrm]
command=/opt/xilinx/xrm/bin/xrmd
environment=LD_LIBRARY_PATH=/opt/xilinx/xrt/lib
priority=500
